# installation/systemd/coziyoo-agent.service
[Unit]
Description=Coziyoo Agent Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/coziyoo/agent-python
EnvironmentFile=/opt/coziyoo/agent-python/.env
ExecStart=/opt/coziyoo/agent-python/.venv/bin/python -m coziyoo_agent.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
